# razer-coolingpad.service
[Unit]
Description=Razer Cooling Pad Fan Curve Controller
After=multi-user.target

[Service]
Type=simple
# Adjust paths below to match your installation
ExecStart=/usr/bin/python3 /opt/razer-coolingpad-linux/razer-coolingpad-fancurve.py -c /opt/razer-coolingpad-linux/curves/balanced.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
